# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create hyphanet user and group
        if ! getent group hyphanet >/dev/null; then
            addgroup --system hyphanet
        fi
        if ! getent passwd hyphanet >/dev/null; then
            adduser --system --ingroup hyphanet --home /var/lib/hyphanet \
                --no-create-home --gecos "Hyphanet daemon" hyphanet
        fi

        # Set ownership of data directory
        chown -R hyphanet:hyphanet /var/lib/hyphanet
        chmod 750 /var/lib/hyphanet

        echo ""
        echo "Hyphanet has been installed."
        echo ""
        echo "To start Hyphanet:"
        echo "  sudo systemctl start hyphanet"
        echo ""
        echo "To enable Hyphanet at boot:"
        echo "  sudo systemctl enable hyphanet"
        echo ""
        echo "Access the web interface at: http://127.0.0.1:8888/"
        echo ""
        ;;
esac

#DEBHELPER#

exit 0
